# rpi-image/systemd/calvin-backend.service
[Unit]
Description=Calvin Dashboard Backend
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=calvin
Group=calvin
WorkingDirectory=/home/calvin/calvin/backend
Environment="PATH=/home/calvin/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DATABASE_URL=sqlite:///home/calvin/calvin/backend/data/db/calvin.db"
Environment="IMAGE_DIR=/home/calvin/calvin/backend/data/images"
# Use venv if it exists (pip installation), otherwise use UV
ExecStart=/bin/bash -c 'cd /home/calvin/calvin/backend && if [ -f .venv/bin/activate ]; then source .venv/bin/activate && uvicorn app.main:app --host 0.0.0.0 --port 8000; else export PATH="/home/calvin/.local/bin:/home/calvin/.cargo/bin:$PATH" && uv run uvicorn app.main:app --host 0.0.0.0 --port 8000; fi'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=calvin-backend

# Security settings
# Note: NoNewPrivileges=true prevents sudo from working, but we need it for reboot
# We'll try systemctl reboot first (which might work without sudo)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Allow writes to Calvin directory (needed for database and data)
ReadWritePaths=/home/calvin/calvin/backend/data /home/calvin/calvin/backend/logs /home/calvin/calvin/backend
# Allow reboot capability (needed for system reboot)
CapabilityBoundingSet=CAP_SYS_BOOT
AmbientCapabilities=CAP_SYS_BOOT

[Install]
WantedBy=multi-user.target

