# rpi-image/systemd/calvin-frontend.service
[Unit]
Description=Calvin Dashboard Frontend (Kiosk Mode)
After=network.target calvin-backend.service
Requires=calvin-backend.service
Wants=graphical.target
# Don't fail if X isn't available yet - it will start via .bash_profile on tty1

[Service]
Type=simple
User=calvin
Group=calvin
Environment="DISPLAY=:0"
Environment="HOME=/home/calvin"
Environment="XAUTHORITY=/home/calvin/.Xauthority"
# Wait for backend to be ready (with timeout)
ExecStartPre=/bin/bash -c 'timeout=60; elapsed=0; until curl -s http://localhost:8000/api/health > /dev/null 2>&1; do if [ $elapsed -ge $timeout ]; then echo "Backend not ready after $timeout seconds"; exit 1; fi; sleep 1; elapsed=$((elapsed+1)); done'
# Wait for X server to be available (with timeout to prevent infinite wait)
ExecStartPre=/bin/bash -c 'timeout=60; elapsed=0; until [ -S /tmp/.X11-unix/X0 ] || pgrep -x Xorg > /dev/null 2>&1; do if [ $elapsed -ge $timeout ]; then echo "X server not available after $timeout seconds"; exit 1; fi; sleep 1; elapsed=$((elapsed+1)); done'
# Small delay to ensure X is fully ready
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/chromium \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --autoplay-policy=no-user-gesture-required \
    --disable-features=TranslateUI \
    --disable-ipc-flooding-protection \
    --disable-background-networking \
    --disable-background-timer-throttling \
    --disable-renderer-backgrounding \
    --disable-backgrounding-occluded-windows \
    --disable-breakpad \
    --disable-component-update \
    --disable-domain-reliability \
    --disable-features=AudioServiceOutOfProcess \
    --disable-hang-monitor \
    --disable-prompt-on-repost \
    --disable-sync \
    --disable-web-resources \
    --metrics-recording-only \
    --mute-audio \
    --no-first-run \
    --no-default-browser-check \
    --no-pings \
    --no-zygote \
    --use-mock-keychain \
    http://localhost:8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=calvin-frontend

[Install]
WantedBy=graphical.target

