# rpi-image/systemd/calvin-backend.service
[Unit]
Description=Calvin Dashboard Backend
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=calvin
Group=calvin
WorkingDirectory=/home/calvin/calvin/backend
Environment="PATH=/home/calvin/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DATABASE_URL=sqlite:///home/calvin/calvin/backend/data/db/calvin.db"
Environment="IMAGE_DIR=/home/calvin/calvin/backend/data/images"
# Use venv if it exists (pip installation), otherwise use UV
# For dev mode, use --reload flag for hot reload
# Check if .dev file exists to enable dev mode
ExecStart=/bin/bash -c 'cd /home/calvin/calvin/backend && if [ -f .venv/bin/activate ]; then source .venv/bin/activate && if [ -f .dev ]; then uvicorn app.main:app --host 0.0.0.0 --port 8000 --reload; else uvicorn app.main:app --host 0.0.0.0 --port 8000; fi; else export PATH="/home/calvin/.local/bin:/home/calvin/.cargo/bin:$PATH" && if [ -f .dev ]; then uv run uvicorn app.main:app --host 0.0.0.0 --port 8000 --reload; else uv run uvicorn app.main:app --host 0.0.0.0 --port 8000; fi; fi'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=calvin-backend

# Security settings
# Note: NoNewPrivileges=true prevents sudo from working
# We'll use systemctl reboot which should work if user has proper permissions
# Or use a helper script with setuid/sudo permissions
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Allow writes to Calvin directory (needed for database and data)
ReadWritePaths=/home/calvin/calvin/backend/data /home/calvin/calvin/backend/logs /home/calvin/calvin/backend

[Install]
WantedBy=multi-user.target

