# rpi-image/systemd/calvin-update.timer
[Unit]
Description=Calvin Dashboard Auto-Update Timer
Requires=calvin-update.service

[Timer]
# Run every hour (configurable via OnUnitActiveSec)
# OnCalendar format doesn't support arbitrary seconds, so we use OnUnitActiveSec
# Changed from 5 minutes to 1 hour to avoid unnecessary rebuilds
OnUnitActiveSec=1h
# Also run 1 hour after boot
OnBootSec=1h
# Run immediately if missed
Persistent=true

[Install]
WantedBy=timers.target

